# supervisord.conf
[supervisord]
nodaemon=false

[program:docker]
command=sh /usr/local/bin/dockerd-entrypoint.sh
autostart=true
autorestart=true
process_name=docker_dind